# files/desktop/config/usr/lib/systemd/user/dunst.service.d/install.conf
[Unit]
After=wayland-session@hyprland.desktop.target

[Install]
WantedBy=graphical-session.target
WantedBy=wayland-session@hyprland.desktop.target
