# files/desktop/config/usr/lib/systemd/user/waybar.service.d/install.conf
[Unit]
After=wayland-session@hyprland.desktop.target

[Install]
WantedBy=wayland-session@hyprland.desktop.target
