# files/desktop/config/usr/lib/systemd/user/waybar.service.d/install.conf


[Install]
WantedBy=wayland-session@hyprland.desktop.target
